# gateway/android-vacuum-gateway.service
[Unit]
Description=Deterministic Xiaomi Android vacuum gateway
After=network-online.target android-mcp-forward.service
Wants=network-online.target

[Service]
Type=simple
User=android-vacuum
Group=android-vacuum
WorkingDirectory=/opt/android-vacuum-gateway
Environment=HOME=/var/lib/android-vacuum-gateway
Environment=PYTHONDONTWRITEBYTECODE=1
UMask=0077
ExecStart=/usr/bin/python3 /opt/android-vacuum-gateway/android_vacuum_gateway.py --config /etc/android-vacuum-gateway.json
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
StateDirectory=android-vacuum-gateway
StateDirectoryMode=0700
ReadWritePaths=/var/lib/android-vacuum-gateway

[Install]
WantedBy=multi-user.target
